# playbooks/roles/DEBIAN12-CIS/templates/etc/profile.d/CIS.sh.j2
# {{ ansible_managed }}

umask {{ deb12cis_bash_umask }}
